# raspberrypi/udev.rules
ACTION=="add" \
, RUN+="/home/pi/udev.sh"
